# ansible/roles/web/templates/gunicorn.socket.j2
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn-{{ ssp_name }}/socket

[Install]
WantedBy=sockets.target